# env/roles/redis/templates/redis-supervisor.conf
[program:redis]
command=redis-server

stdout_logfile={{ log_dir }}/redis.log
stderr_logfile={{ log_dir }}/redis.log