# env/roles/socket/templates/socket-supervisor.conf
[program:socket]
command=/usr/bin/python {{ code_dir }}/socket/app.py
stdout_logfile={{ log_dir }}/socket.log
stderr_logfile={{ log_dir }}/socket.log
autorestart=true
