# env/roles/socket/templates/socket-supervisor.conf
[program:socket]
command=python -u app.py
directory=/exploder/socket
stdout_logfile={{ log_dir }}/socket.log
stderr_logfile={{ log_dir }}/socket.log
autorestart=true
user=root
