# env/roles/socket/templates/socket-supervisor.conf
[program:socket]
command=python /exploder/socket/app.py
stdout_logfile={{ log_dir }}/socket.log
stderr_logfile={{ log_dir }}/socket.log
autorestart=true
user=root
